# _apps/pik-lk-notification-backend-prod-ruby/runit/log/run
#!/bin/sh
exec svlogd -tt /var/log/runit/{{ app_name }}
